# centos/7.4/aws/DCOS-1.13.3/docker-18.09.2-ce/selinux_disabled/mount_tmp_noexec.service
[Unit]
Description=Mount tmp noexec,nosuid,nodev at startup

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/sbin/mount_tmp_noexec.sh

[Install]
WantedBy=multi-user.target
